# ansible/files/systemd/agrocomm.service
[Unit]
Description=AgroComm Next.js Application
After=network.target

[Service]
Type=simple
User=nginx
MemoryAccounting=true
MemoryHigh=768M
MemoryMax=1024M
Environment=NODE_OPTIONS=--max-old-space-size=512
WorkingDirectory=/var/www/agrocomm
ExecStart=/home/nginx/.local/share/pnpm/pnpm run start
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10s
# Limitar tentativas de reinicialização para evitar falhas em cascata
StartLimitIntervalSec=60
StartLimitBurst=3
# Garantir que o serviço seja considerado "falho" se não iniciar corretamente
TimeoutStartSec=90s
# Verificação de saúde: tenta até 6x com 5s de intervalo (máx ~30s)
ExecStartPost=/bin/sh -c 'for i in 1 2 3 4 5 6; do sleep 5 && curl -s -f http://localhost:4000/api/health > /dev/null 2>&1 && exit 0; done; exit 1'
# Arquivo de ambiente
EnvironmentFile=-/var/www/agrocomm/.env.production
# Configurações de segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target