# ansible/files/systemd/agrocomm-bulletins.timer
[Unit]
Description=AgroComm Bulletins Timer - Boletins e alertas de assinatura
Requires=agrocomm-bulletins.service

[Timer]
# Roda de hora em hora; a lógica de agenda (dias/horários) está no banco de dados
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
